# test/unit/data/test_quadlet/oci_port/oci-model-port.container
[Unit]
Description=RamaLama oci-model-port AI Model Service
After=local-fs.target

[Container]
AddDevice=-/dev/accel
AddDevice=-/dev/dri
AddDevice=-/dev/kfd
AddDevice=nvidia.com/gpu=all
Image=testimage
RunInit=true
Environment=HOME=/tmp
Exec=
SecurityLabelDisable=true
DropCapability=all
NoNewPrivileges=true
PublishPort=0.0.0.0:8080:8080
Mount=type=image,source=registry.example.com/model:latest,destination=/mnt/models,subpath=/models,readwrite=false

[Install]
WantedBy=multi-user.target default.target

